# deploy/vai.service
[Unit]
Description=VAI local UI (Hari random mode)
After=network.target

[Service]
Type=simple
User=vai
Group=vai
WorkingDirectory=/opt/vai
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/vai/.venv/bin/python -m vai serve --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target